# scripts/ops/market-prices.service
[Unit]
Description=EDB Market Prices Refresh
After=network-online.target
Requires=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/projects/EVEDataBrowser
EnvironmentFile=/etc/edb/market.env
ExecStart=/usr/bin/env npm --prefix app/api run ingest:market:prices
Nice=5
StandardOutput=journal
StandardError=inherit

[Install]
WantedBy=multi-user.target
